# kak-lsp.service
[Unit]
Description=Kak LSP daemon

[Service]
ExecStart=/bin/sh -c 'kak-lsp'
Type=simple
Restart=always

[Install]
WantedBy=default.target
