# kak-lsp.service
[Unit]
Description=Kak LSP daemon

[Service]
ExecStart=/home/lars/.local/bin/kak-lsp
Type=simple
Restart=always

[Install]
WantedBy=default.target
